# linux/systemd/brain-dashboard.service.template
[Unit]
Description=Brain AI Hub — Dashboard (FastAPI)
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=@USER@
Group=@USER@
WorkingDirectory=@REPO@/dashboard
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONIOENCODING=utf-8
Environment=OLLAMA_HOST=127.0.0.1:@OLLAMA_PORT@
Environment=BRAIN_DATA_DIR=@DATA@
ExecStart=@VENV@/bin/uvicorn app:app --host @BIND@ --port @DASH_PORT@
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
