# linux/systemd/brain-mcp-gateway.service.template
[Unit]
Description=Brain AI Hub — MCP HTTP/SSE Gateway (mcp-proxy, internal only)
After=network.target brain-dashboard.service
Wants=brain-dashboard.service

[Service]
Type=simple
User=@USER@
Group=@USER@
WorkingDirectory=@REPO@/dashboard
Environment=PYTHONUNBUFFERED=1
Environment=OLLAMA_HOST=127.0.0.1:11434
# mcp-proxy (PyPI: mcp-proxy, sparfenyuk) multiplexes brain-vault/brain-library
# (filesystem MCP servers, via the npm-installed mcp-server-filesystem binary)
# and brain-rag (mcp_rag.py, unnamed/default server) behind one HTTP/SSE port.
# Binds to loopback-only — public clients hit brain-mcp-auth-proxy on the
# configured public port, which Bearer-checks and reverse-proxies here.
ExecStart=@VENV@/bin/mcp-proxy --host 127.0.0.1 --port @MCP_UPSTREAM_PORT@ --pass-environment \
  --named-server brain-vault "mcp-server-filesystem @DATA@/vault" \
  --named-server brain-library "mcp-server-filesystem @DATA@/library" \
  @VENV@/bin/python mcp_rag.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
