# linux/systemd/brain-mcp-gateway.service.template
[Unit]
Description=Brain AI Hub — MCP SSE Gateway (supergateway, internal only)
After=network.target brain-dashboard.service
Wants=brain-dashboard.service

[Service]
Type=simple
User=@USER@
Group=@USER@
WorkingDirectory=@REPO@
Environment=NODE_ENV=production
Environment=OLLAMA_HOST=127.0.0.1:11434
# IMPORTANT: supergateway listens on @MCP_UPSTREAM_PORT@ (loopback-only).
# Public clients hit brain-mcp-auth-proxy on @MCP_PORT@, which Bearer-checks
# and reverse-proxies here. If supergateway does not honor --host on your
# version, ensure firewall blocks @MCP_UPSTREAM_PORT@ from LAN:
#   sudo ufw deny @MCP_UPSTREAM_PORT@/tcp
ExecStart=/usr/bin/npx -y supergateway --multiServerConfig @REPO@/linux/mcp-multi-server.json --port @MCP_UPSTREAM_PORT@ --host 127.0.0.1 --outputTransport sse --ssePath /sse --messagePath /message
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
